# docker/crontabs/root
* * * * * /sbin/sudo -u www-data /app/bin/console repman:proxy:sync-releases
* * * * * /sbin/sudo -u www-data /app/bin/console repman:proxy:sync-metadata
0 */2 * * * /sbin/sudo -u www-data /app/bin/console repman:security:update-db
0 0 * * * /sbin/sudo -u www-data /app/bin/console repman:telemetry:send
0 1 * * * /sbin/sudo -u www-data /app/bin/console repman:package:clear-old-dists
